# packaging/analysisview.cron
MAILTO=""
*/3 * * * * apache analysisview-update > /dev/null
