# packaging/analysiscrab2.cron
MAILTO=
*/5 * * * * apache analysisviewcrab2-update
