# packaging/analysiscrab2.cron
MAILTO=""
*/5 * * * * root analysisviewcrab2-update > /dev/null
