# packaging/analysisview.cron
MAILTO=""
*/3 * * * * root analysisview-update > /dev/null
