# packaging/cleanup.cron
MAILTO=""
* */4 * * * root sh /etc/prodview/cleanup.sh > /dev/null
