# packaging/cmsconnectview.cron
MAILTO=""
*/3 * * * * root cmsconnect-update > /dev/null
