# packaging/condormonitoringcron.cron
MAILTO=""
*/15 * * * * root /data/prodview/src/CondorMonitoring/JobCounter.sh &> /dev/null
