# packaging/history.cron
MAILTO=""
*/10 * * * * root history-update -t production > /dev/null
*/10 * * * * root history-update -t analysis > /dev/null
