# packaging/history.cron
MAILTO=""
*/10 * * * * root history-update > /dev/null
