# packaging/institutionalview.cron
MAILTO=""
*/3 * * * * root institutional-update > /dev/null
