# packaging/prodview.cron
MAILTO=""
*/3 * * * * apache prodview-update > /dev/null
