# packaging/prodview.cron
MAILTO=""
*/3 * * * * root prodview-update > /dev/null
