# packaging/totalview.cron
MAILTO=""
*/3 * * * * apache totalview-update > /dev/null
