# packaging/totalview.cron
MAILTO=""
*/3 * * * * root totalview-update > /dev/null
