# packaging/utilization.cron
MAILTO=""
10 */6 * * * root prodview-update > /dev/null
