# packaging/utilization.cron
MAILTO=""
10 */1 * * * root utilization > /dev/null
