# packaging/viewcleanup.cron
MAILTO=
11 1,13 * * * apache /etc/prodview/cleanup.sh
